# app-blocker-daemon.service
[Unit]
Description=App Blocker Daemon
After=network.target

[Service]
Type=simple
ExecStart=PYTHON_PATH_PLACEHOLDER REPO_PATH_PLACEHOLDER/daemon.py
WorkingDirectory=REPO_PATH_PLACEHOLDER
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target